# services/fotoroulette_processing.service
[Unit]
Description=Fotoroulette processing layer

[Service]
Type=simple
ExecStart=/home/thomas/PycharmProjects/raspberry-pi/processing_main.py

[Install]
WantedBy=multi-user.target